# .config/autostart/kpie.desktop
[Desktop Entry]
Name=KPie
Comment=Window Utility
TryExec=/usr/bin/kpie
StartupNotify=false
Terminal=false
Type=Application
